# .swoole-supervisor.conf
[program:swoole-develop]
directory=/var/www/swoole-develop/laravel
command=php bin/laravels start -i
autostart=true
autorestart=true
user=root
stdout_logfile=/var/log/supervisor/swoole-develop.log